# deploy/vibebridge-tunnel.service
[Unit]
Description=AI Project Lab Tunnel (Auto-switching)
After=vibebridge.service
Requires=vibebridge.service

[Service]
Type=simple
User=akliedrak
WorkingDirectory=/home/user/workspace/vibebridge
Environment=PATH=/home/user/workspace/vibebridge/.venv/bin:/home/user/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/user/workspace/vibebridge/scripts/install/auto_recovery.sh daemon
Restart=always
RestartSec=10
StartLimitIntervalSec=0
# 资源限制
MemoryMax=300M
MemorySwapMax=100M
LimitNOFILE=65536
LimitNPROC=65536
# 日志输出
StandardOutput=append:/home/user/workspace/vibebridge/logs/tunnel.log
StandardError=append:/home/user/workspace/vibebridge/logs/tunnel.log
# 安全选项
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/user/workspace/vibebridge/logs

[Install]
WantedBy=multi-user.target