# deploy/vibebridge.service
[Unit]
Description=AI Project Lab Server
After=network.target

[Service]
Type=simple
User=akliedrak
WorkingDirectory=/home/user/workspace/vibebridge
Environment=PATH=/home/user/workspace/vibebridge/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/user/workspace/vibebridge/.venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StartLimitIntervalSec=0
# 资源限制
MemoryMax=500M
MemorySwapMax=200M
LimitNOFILE=65536
LimitNPROC=65536
# 日志输出
StandardOutput=append:/home/user/workspace/vibebridge/logs/server.log
StandardError=append:/home/user/workspace/vibebridge/logs/server.log
# 安全选项
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/user/workspace/vibebridge/logs /home/user/workspace/vibebridge/data

[Install]
WantedBy=multi-user.target